# Shell/checkvpn.service
[Unit]
Description = VPN serivce check
After = network.target

[Service]
ExecStart = /usr/sbin/checkvpn.sh

[Install]
WantedBy = multi-user.target